# canvas_grit_sync.timer
[Unit]
Description=Run Canvas Grit sync every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
